# kumonoboru.service.j2
[Unit]
Description=Kumonoboru - cloud backup utility

[Service]
ExecStart=/usr/bin/kumonoboru

[Install]
WantedBy=multi-user.target
