# kumonoboru.service.j2
[Unit]
Description=Kumonoboru - cloud backup utility

[Service]
ExecStart=/bin/bash /usr/local/bin/kumonoboru.sh

[Install]
WantedBy=multi-user.target
