# kumonoboru.service
[Unit]
Description=Kumonoboru - Restic backup with Prometheus reporting
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/kumonoboru/env
ExecStart=/usr/bin/kumonoboru

[Install]
WantedBy=multi-user.target
