# kumonoboru-prune.service.j2
[Unit]
Description=Kumonoboru - Cloud Backup Utility

[Service]
ExecStart=/usr/bin/kumonoboru --clean

[Install]
WantedBy=multi-user.target
